# media-compressor.service
[Unit]
Description=Media Compressor Job
After=network.target

[Service]
Type=oneshots
ExecStart=/usr/bin/python3 /home/devuser/main.py /mnt/media
User=devuser